# systemd/dr-audio-watch.service
[Unit]
Description=resolve-audio-fix — DaVinci Resolve AAC Audio Watcher
Documentation=https://github.com/owlivion/resolve-audio-fix
After=default.target

[Service]
Type=simple
ExecStart=%h/.local/bin/davinci-audio-watchd
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
